# systemd/singbox-rule-ui.service
[Unit]
Description=Sing-box custom rule web UI
After=network-online.target sing-box.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/singbox-rule-ui
Environment=RULE_UI_HOST=0.0.0.0
Environment=RULE_UI_PORT=8088
ExecStart=/usr/bin/python3 /opt/singbox-rule-ui/app.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
