# systemd/update-sing-box-rules-jsdelivr.service
[Unit]
Description=Update sing-box geosite/geoip rule sets via jsDelivr
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-sing-box-rules-jsdelivr

